# output/files/etc/hotplug.d/module/20-proslic
#!/bin/sh
# SPDX-License-Identifier: GPL-2.0-or-later
# Hotplug script: create /dev/proslic{0,1} when slic_si3226x loads.
#
# The slic-si3226x driver registers misc devices automatically, but this
# script ensures correct permissions and group ownership on the device
# nodes, and handles any udev/mdev gaps on minimal OpenWrt installs.

[ "$MODALIAS" = "slic_si3226x" ] || [ "$DEVPATH" != "${DEVPATH%slic_si3226x*}" ] || exit 0

case "$ACTION" in
add)
	# Wait briefly for misc device registration
	tries=0
	while [ $tries -lt 20 ] && [ ! -e /dev/proslic0 ]; do
		sleep 0.25
		tries=$((tries + 1))
	done

	# Ensure device nodes exist (fallback: create from /proc/devices)
	if [ ! -e /dev/proslic0 ]; then
		major=$(awk '/proslic/{print $1}' /proc/devices 2>/dev/null)
		if [ -n "$major" ]; then
			mknod /dev/proslic0 c "$major" 0 2>/dev/null
			mknod /dev/proslic1 c "$major" 1 2>/dev/null
		fi
	fi

	# Set permissions: 660, group dialout
	for dev in /dev/proslic0 /dev/proslic1; do
		[ -e "$dev" ] || continue
		chmod 660 "$dev"
		chgrp dialout "$dev" 2>/dev/null || true
	done

	logger -t proslic-hotplug "SLIC devices configured"
	;;
remove)
	rm -f /dev/proslic0 /dev/proslic1
	logger -t proslic-hotplug "SLIC devices removed"
	;;
esac
